# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid.cron
#
MAILTO=john.florian@dart.biz
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# The main mirrmaid.conf is used for:
#       1) gathering PICAPS definitions for local backup
#       2) YUM repository distribution to mirrors local at each plant
2,17,32,47 * * * * mirrmaid     nice ionice -c3 mirrmaid -c /etc/mirrmaid/mirrmaid.conf
