# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid-mariadb.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid-mariadb.cron
#
MAILTO=john.florian@dart.biz
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# All jobs declared here will use this proxy configuration.
RSYNC_PROXY=mdct-proxy:8080

# Fedora mirroring (rawhide anyway).
15 * * * * mirrmaid           nice ionice -c3 mirrmaid -c /etc/mirrmaid/mirrmaid-mariadb.conf
