# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid-yum-fanout.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# files/private/mdct-00fs.dartcontainer.com/mirrmaid/mirrmaid-yum-fanout.cron
#
MAILTO=john.florian@dart.biz
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# YUM repository distribution to mirrors local at each plant
29,59 * * * * mirrmaid     nice ionice -c3 mirrmaid -c /etc/mirrmaid/mirrmaid-yum-fanout.conf
