# environments/production/modules/dart/files/puppet/tools/puppet-report-cleanup.cron
# vim:filetype=crontab
#
# Synopsis:
#       Purge oldest puppet reports and compress all but the most recent of
#       those that remain.
#
# Notes:
#       This file is managed by puppet.

PATH="/sbin:/bin:/usr/sbin:/usr/bin"
17 5 * * * puppet       nice ionice -c 3  puppet-report-cleanup
