# environments/production/modules/dart/templates/mirrmaid/mirrmaid-mariadb.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# Origin: modules/dart/templates/mirrmaid/mirrmaid-mariadb.cron
#
MAILTO=john.florian@dart.biz
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# All jobs declared here will use this proxy configuration.
RSYNC_PROXY=<%= @rsync_proxy %>

# MariaDB mirroring
*/30 * * * * mirrmaid       nice ionice -c3 mirrmaid -c /etc/mirrmaid/mirrmaid-mariadb.conf
