# files/private/mdct-00fs.dartcontainer.com/picaps/picaps-backup-agent.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# cron job for picaps-backup-agent

SOURCES="mdct-01pi mdct-02pi mdct-03pi mdct-04pi mdct-05pi mdct-10pi mdct-11pi mdct-15pi mdct-19pi mdct-25pi mdct-39pi mdct-47pi mdct-52pi mdct-55pi mdct-72pi mdct-73pi mdct-75pi mdct-76pi mdct-77pi mdct-78pi mdct-95pi mdct-98pi mdct-0302pi mdct-0310pi mdct-0314pi"
SOURCE_RSYNC_URL="/backup"
TARGET_PATH="/storage/remote-backups/picaps/event-sets"

# Maintain the off-site event set archives.
# NB: Definition archives are managed by mirrmaid.
39 * * * *  root        nice ionice -c3 picaps-backup-agent $SOURCE_RSYNC_URL $TARGET_PATH $SOURCES

# Produce and deliver a report on the status of the off-site definition data mirrors.
18 5 * * Thu root       nice ionice -c3 verify-picaps-definitions-backup
