# modules/dart/templates/mirrmaid/mirrmaid-rpmfusion.cron
# ***** THIS FILE IS MANAGED BY PUPPET *****
# Origin: modules/dart/templates/mirrmaid/mirrmaid-rpmfusion.cron
#
MAILTO=john.florian@dart.biz
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# All jobs declared here will use this proxy configuration.
RSYNC_PROXY=<%= @rsync_proxy %>

# RPM Fusion mirroring
*/18 * * * * mirrmaid       nice ionice -c3 mirrmaid -c /etc/mirrmaid/mirrmaid-rpmfusion.conf
