# files/private/mdct-02ut.dartcontainer.com/xinetd.d/rsync
service rsync
{
	disable	        = no
	flags		= IPv6
	socket_type     = stream
	wait            = no
	user            = root
	server          = /usr/bin/rsync
	server_args     = --daemon
	log_on_failure  += USERID

        # Avoid thundering herds, better to serve a few a time through
        # completion than to poorly serve all.
        max_load        = 3
        nice            = 10
        per_source      = 1
}
